# crontabs/powerline-daemmon.crontab
*/10 * * * * /bin/ps powerline-daemonn -o pid= >/dev/null || /usr/bin/powerline-daemon -q
