# crontabs/update-cron.sh
#!/bin/bash
cat ${HOME}/crontabs/*.crontab | /usr/bin/crontab -
